# roles/consul/files/consul.conf
description "consul"
author "Ian Unruh <ianunruh@gmail.com>"

start on virtual-filesystems
stop on runlevel [016]

respawn

setuid consul
setgid consul

exec consul agent -config-dir=/etc/consul/conf.d >> /var/log/consul/agent.log 2>&1
